# conf/crontab
# установка файла: sudo crontab -u autocatalogs /srv/autocatalogs/conf/crontab

HOME=/srv/autocatalogs
SERVER_ADMIN=""

# Обновление справочников
0 8 * * * npm run catalogs:update > /dev/null
